# build/supervisor/beat.conf
[program:beats]
directory=/opt/webapp
command=/usr/local/bin/celery -A tasks beat -l debug --scheduler django_celery_beat.schedulers:DatabaseScheduler
autostart=true
autorestart=true
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0